# profileinit/etc/rc.local
#!/bin/sh

echo 'Starting local daemons'
