# profileinit/etc/rc.local
#!/bin/sh

log 'Starting local daemons'

:

log 'Local daemons started'
