# profileinit/profileinit.service
[Unit]
Description=profileinit service

[Service]
ExecStart=/usr/local/bin/profileinit

[Install]
WantedBy=multi-user.target
